# Non Docker Setup/systemd/fbreelz_refresh.service
# ## version 1
[Unit]
Description=FBReelz - Refresh saved items and download cache
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/fbreelz
Environment=TZ=Europe/London
ExecStart=/bin/bash -lc 'source /opt/fbreelz/.venv/bin/activate && python /opt/fbreelz/fbreelz_phase1_playwright.py --max 30 && python /opt/fbreelz/fbreelz_phase2_resolve.py --download && python /opt/fbreelz/make_cache_playlist.py --base-url http://192.168.1.226/cache/ --out /opt/fbreelz/data/fbreelz_cache_http.m3u'
